# provisioning/systemd/sis-screentime.service
[Unit]
Description=Enforce screen time restrictions

[Service]
Type=oneshot
ExecStart=/usr/local/sis/enforce-screentime.sh
