# provisioning/systemd/sis-remote-wipe.path
[Unit]
Description=Watch for SIS remote wipe flag

[Path]
PathExists=/etc/sis/REMOTE_WIPE

[Install]
WantedBy=multi-user.target
