# provisioning/systemd/sis-mdm-telemetry.service
[Unit]
Description=SIS MDM telemetry upload
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sis/mdm-agent.sh telemetry
