# provisioning/systemd/sis-zerotouch.target
[Unit]
Description=SIS Zero-touch Setup Target
Requires=sis-zerotouch.service
