# provisioning/systemd/sis-remote-wipe.service
[Unit]
Description=Execute SIS remote wipe

[Service]
Type=oneshot
ExecStart=/usr/local/sis/remote-wipe.sh
